# api/supervisord.conf
[supervisord]
nodaemon=true

[program:postgres]
command=/usr/lib/postgresql/14/bin/postgres -D /var/lib/postgresql/data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/postgres.log

[program:redis]
command=redis-server
redirect_stderr=true
stdout_logfile=/var/log/supervisor/redis.log

[program:web]
command=gunicorn -b 0.0.0.0:8000 app:app
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/web.log

[program:celery-worker]
command=celery -A celery_app.celery worker --loglevel=info
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/celery_worker.log

[program:celery-beat]
command=celery -A celery_app.celery beat --loglevel=info
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/celery_beat.log

